# rootfs/etc/services.d/swift/run
#!/bin/sh
swift-init main start no-daemon
