# rootfs/etc/services.d/confd/run
#!/bin/sh
confd -watch -backend etcd 2>&1
